# service/log/run
#!/bin/sh

exec setuidgid takatoshi multilog t ./main

